# noatime.service
[Unit]
After=graphical.target

[Service]
ExecStart=mount -o remount,rw,noatime /

[Install]
WantedBy=graphical.target
